# tomcat-mysql/supervisord.ini
[supervisord]
nodaemon=true

[program:mysql]
command=/mysql-start.sh
numprocs=1
autostart=true
autorestart=true
user=root

[program:tomcat]
command=/usr/local/tomcat/bin/catalina.sh run
numprocs=1
autostart=true
autorestart=true
user=root